# launchd/homebrew.mopidy.mopidy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>homebrew.mopidy.mopidy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/opt/mopidy/bin/mopidy</string>
    <string>--config</string>
    <string>/opt/tunebot/mopidy.conf</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>